# supervisorconf/ftp.conf
[program:radius]
command=/usr/sbin/radiusd -d /etc/raddb -xX
priority=1
redirect_stderr=true

