# supervisorconf/portal.conf
[program:protocal]
command=/usr/share/adfi/protocal/syler -config="/usr/share/adfi/protocal/sylerconf/syler.conf"
priority=1
redirect_stderr=true

