# supervisorconf/sms.conf
[program:sms-adapter]
command=/usr/share/adfi/smsadapter/smsadapter -config=/usr/share/adfi/smsadapter/smsadapter.conf
priority=1
redirect_stderr=true

